# templates/default/collectd.upstart.conf.erb
description "collectd"

start on runlevel [2345]
stop on runlevel [s016]
respawn
expect fork
respawn limit 10 5

script
    exec collectd -C /etc/collectd/collectd.conf 
end script
